# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true

[program:node]
command=npm start
directory=/usr/src/app
autostart=true
autorestart=true